# crontab
0 */2 * * * php /var/www/Keepstar/cron.php
